# networkmanager/AutoWifiConnect.timer
[Unit]
Description=Timer for periodically scanning for WiFi.

[Timer]
OnBootSec=120s
OnUnitActiveSec=120s

[Install]
WantedBy=timers.target