# networkmanager/AutoWifiConnect.timer
[Unit]
Description=Timer for periodically scanning for WiFi.

[Timer]
OnBootSec=60s
OnUnitActiveSec=60s

[Install]
WantedBy=timers.target