# networkmanager/AutoWifiConnect.timer
[Unit]
Description=Timer for periodically scanning for WiFi.

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target